# deprecated/quadlet/ddclient/ddclient.container
[Container]
ContainerName=ddclient
Image=lscr.io/linuxserver/ddclient:latest
Volume=ddclient-config.volume:/config:z

[Service]
Restart=always