# quadlet/bots/bots.target
[Unit]
Description=Bots
After=net.target
Wants=net.target

[Install]
WantedBy=multi-user.target default.target