# quadlet/smarthome/smarthome.target
[Unit]
Description=Smarthome
After=net.target
Wants=net.target